# scripts/templates/memory-short-rest.timer
[Unit]
Description=Run short-rest maintenance daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
